# systemd/chkpt.service
[Unit]
Description=Dump service, periodical

[Service]
Type=simple
ExecStart=/bin/bash -c "$$HOME/bin/dumpme"

[Install]
WantedBy=default.target
